# example/wakelet.service
[Unit]
Description=Wakelet — HomeKit IoT bridge
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/wakelet
ExecStartPre=/bin/mkdir -p /var/lib/wakelet
ExecStartPre=/bin/mkdir -p /etc/wakelet/private
ExecStart=/opt/wakelet/.venv/bin/python driver.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
