# example/wakelet.service
[Unit]
Description=Wakelet — HomeKit IoT bridge
After=network.target

[Service]
Type=simple
User=wakelet
Group=wakelet
WorkingDirectory=/home/wakelet/wakelet
ExecStart=/home/wakelet/wakelet/.venv/bin/python driver.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
